# luci-app-network-detector/root/etc/uci-defaults/99-luci-network-detector
#!/bin/sh
# LuCI 网络检测器 - 首次安装初始化
# OpenWrt 首次启动时自动执行，执行后自删除

[ -n "${IPKG_INSTROOT}" ] && exit 0

# =============================================================
# 自动检测 OpenClash API Secret
# =============================================================
detect_openclash_secret() {
    local secret=""

    if ! command -v uci > /dev/null 2>&1; then
        return 1
    fi

    # 1) 尝试从 /etc/config/openclash 的 UCI 读取
    secret=$(uci -q get openclash.@config[0].dashboard_password 2>/dev/null)
    if [ -n "$secret" ]; then
        echo "network-detector: 从 OpenClash UCI 获取到密钥"
        printf '%s' "$secret"
        return 0
    fi

    # 2) uci show 兜底（某些版本 key 名不同）
    secret=$(uci -q show openclash 2>/dev/null | grep -iE '(password|secret|token)' | head -1 | sed "s/^[^=]*='*//;s/'$//")
    if [ -n "$secret" ]; then
        echo "network-detector: 从 OpenClash UCI show 获取到密钥"
        printf '%s' "$secret"
        return 0
    fi

    # 3) 从 OpenClash YAML 提取（多目录、.yaml/.yml、不依赖行首锚点）
    for dir in /etc/openclash/config /etc/openclash/custom /etc/openclash; do
        [ -d "$dir" ] || continue
        for conf in "$dir"/*.yaml "$dir"/*.yml; do
            [ -f "$conf" ] || continue
            secret=$(grep -hm1 'secret:' "$conf" 2>/dev/null | sed "s/.*secret:[[:space:]]*['\"]\?//;s/['\"].*$//;s/[[:space:]]*$//" | head -1)
            if [ -n "$secret" ]; then
                echo "network-detector: 从 OpenClash YAML ($conf) 获取到密钥"
                printf '%s' "$secret"
                return 0
            fi
        done
    done

    return 1
}

# 尝试自动填入密钥
OPENCLASH_SECRET=$(detect_openclash_secret)
if [ -n "$OPENCLASH_SECRET" ]; then
    uci -q set "network-detector.@clash[0].secret=$OPENCLASH_SECRET"
    uci commit network-detector
fi

# =============================================================
# 自动检测 Clash API 地址
# =============================================================
# 检测已安装的 Clash 控制器的默认 API 端口
detect_clash_api() {
    # OpenClash 默认端口
    if [ -f /etc/openclash/custom/openclash_custom_rules.sh ] 2>/dev/null; then
        # OpenClash 存在，尝试从配置提取
        local api_port=$(uci -q get openclash.@config[0].cn_port 2>/dev/null)
        if [ -n "$api_port" ]; then
            echo "http://127.0.0.1:$api_port"
            return 0
        fi
    fi
    # 默认
    echo "http://127.0.0.1:9090"
}

CLASH_API=$(detect_clash_api)
uci -q set "network-detector.@clash[0].api_url=$CLASH_API"
uci commit network-detector

# =============================================================
# 确保 cron 服务已启用
# =============================================================
/etc/init.d/cron enable 2>/dev/null

# 启动 network-detector 服务（如果已安装核心包）
if [ -x /etc/init.d/network-detector ]; then
    /etc/init.d/network-detector enable 2>/dev/null
    /etc/init.d/network-detector start 2>/dev/null
fi

exit 0
