# network-detector/files/etc/init.d/network-detector
#!/bin/sh /etc/rc.common
# ========================================
# network-detector procd 服务
# 管理 crontab 定时任务
# ========================================

START=99
USE_PROCD=1

PROG="/usr/bin/network-detector"
CRON_MARKER="# --- network-detector ---"

service_triggers() {
    procd_add_reload_trigger "network-detector"
}

start_service() {
    # ===== 主检测服务 =====
    local enabled=$(uci -q get network-detector.@main[0].enabled 2>/dev/null)
    if [ "$enabled" = "0" ]; then
        echo "network-detector 已被禁用，跳过启动"
        return 0
    fi

    # 检查是否有启用的检测任务
    local has_task=0
    local idx=0
    while uci -q get "network-detector.@task[$idx]" > /dev/null 2>&1; do
        local ten=$(uci -q get "network-detector.@task[$idx].enabled" 2>/dev/null)
        if [ "$ten" = "1" ]; then
            has_task=1
            break
        fi
        idx=$((idx + 1))
    done

    if [ "$has_task" = "0" ]; then
        echo "network-detector: 无启用的检测任务，跳过启动 请在 LuCI 中至少添加一个检测任务"
        return 0
    fi

    # 立即运行一次
    [ -x "$PROG" ] && $PROG &

    # 配置 crontab
    setup_cron
}

stop_service() {
    # 移除 crontab 条目
    if [ -f /etc/crontabs/root ]; then
        sed -i "/${CRON_MARKER}/d" /etc/crontabs/root
        /etc/init.d/cron restart 2>/dev/null
    fi
    echo "network-detector 已停止"
}

reload_service() {
    stop_service
    start_service
}

setup_cron() {
    local ival=$(uci -q get network-detector.@main[0].interval_value 2>/dev/null)

    # 兼容旧版 crontab_interval 写法
    if [ -z "$ival" ]; then
        local old_cron=$(uci -q get network-detector.@main[0].crontab_interval 2>/dev/null)
        if [ -n "$old_cron" ]; then
            # 尝试从 cron 表达式反推: */5 → 5, */1 → 1
            case "$old_cron" in
                */1\ *\ *\ *\ *)  ival="1" ;;
                */5\ *\ *\ *\ *)   ival="5" ;;
                */10\ *\ *\ *\ *)  ival="10" ;;
                */15\ *\ *\ *\ *)  ival="15" ;;
                */30\ *\ *\ *\ *)  ival="30" ;;
                *)                ival="5" ;;
            esac
        fi
    fi

    : ${ival:="5"}

    # 校验数值
    case "$ival" in ''|*[!0-9]*) ival="5" ;; esac
    [ "$ival" -lt 1 ] && ival="1"

    local cronexpr="*/${ival} * * * *"

    # 先移除旧条目
    if [ -f /etc/crontabs/root ]; then
        sed -i "/${CRON_MARKER}/d" /etc/crontabs/root
    fi

    # 添加新条目
    echo "${cronexpr} ${PROG} > /dev/null 2>&1 ${CRON_MARKER}" >> /etc/crontabs/root

    # 重启 cron 服务
    /etc/init.d/cron enable 2>/dev/null
    /etc/init.d/cron restart 2>/dev/null

    echo "network-detector crontab 已设置: ${cronexpr} (每 ${ival} 分钟)"
}
